# templates/p4broker_instance_init.erb
#!/bin/bash
#------------------------------------------------------------------------------
# Init script for Helix Broker.
#
# chkconfig: - 90 30
# description: Start Helix Broker for instance REPL_SDP_INSTANCE.
#------------------------------------------------------------------------------

export SDP_INSTANCE=<%= @title %>
/p4/common/bin/p4broker_base $SDP_INSTANCE $@
