# templates/p4p_instance_init.erb
#!/bin/bash
#------------------------------------------------------------------------------
# Init script for Helix Proxy.
#
# chkconfig: - 90 30
# description: Start Helix Proxy for instance REPL_SDP_INSTANCE.
#------------------------------------------------------------------------------

export SDP_INSTANCE=<%= @title %>
/p4/common/bin/p4p_base $SDP_INSTANCE $@

